# clazzvilla-node.conf
[program:classbuddynodeserver]
directory=/classbuddynodeserver
command=serve -l 3000 -s build
stdout_logfile=/classbuddynodeserver/log/stdout.log
