# clazzvilla-node.conf
[program:classbuddynodeserver]
directory=/classbuddynodeserver
command=npm run dev
stdout_logfile=/classbuddynodeserver/log/stdout.log
