# clazzvilla-server.conf
[program:classbuddyserver]
directory=/classbuddyserver
command=java -jar /classbuddyserver/lib/rest-service-0.0.1-SNAPSHOT.jar
stdout_logfile=/classbuddyserver/log/stdout.log
